# cta-tracker/rootfs/etc/services.d/cta-tracker/run
#!/usr/bin/with-contenv bashio

declare api_key
declare poll_interval
declare server_port

api_key=$(bashio::config 'api_key')
poll_interval=$(bashio::config 'poll_interval')
server_port=$(bashio::config 'server_port')

if [ -z "${api_key}" ]; then
    bashio::log.fatal "API key is required. Please set the api_key option."
    exit 1
fi

export CTA_API_KEY="${api_key}"
export POLL_INTERVAL="${poll_interval:-60}"
export SERVER_PORT="${server_port:-8080}"

bashio::log.info "Starting CTA Train Tracker (poll: ${POLL_INTERVAL}s, port: ${SERVER_PORT})"

mkdir -p /share/cta_data

exec /usr/bin/python3 /usr/bin/cta_tracker.py
